# xtrafiles/units/ushutdn.path
[Path]
PathModified=/tmp/ushutdn.tmp

[Install]
WantedBy=multi-user.target

